# deb/dnsminder_0.0.3_all/DEBIAN/postinst
#!/bin/bash

systemctl enable --now resolvconf.service

exit 0
